# sentinel/systemd/sentinel-scanner.timer
[Unit]
Description=Sentinel Port Scanner Timer
Documentation=https://github.com/speckitime/sentinel-firewall

[Timer]
# Run 5 minutes after boot
OnBootSec=5min
# Run daily at 3:00 AM
OnCalendar=*-*-* 03:00:00
# Persist across reboots (run missed timers on boot)
Persistent=true
Unit=sentinel-scanner.service

[Install]
WantedBy=timers.target
